# config/supervisor/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid

[program:openclaw]
command=/usr/local/bin/openclaw-start-guard.sh
directory=/app
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=NODE_ENV="production",HOME="/data",OPENCLAW_HOME="/data",OPENCLAW_STATE_DIR="/data/.openclaw",OPENCLAW_WORKSPACE_DIR="/data/workspace",XDG_CONFIG_HOME="/data/.config"

[program:paperclip-overnight-runner]
command=/bin/bash -lc 'if [ -x /data/workspace/scripts/paperclip-overnight-runner.sh ]; then exec /data/workspace/scripts/paperclip-overnight-runner.sh; else echo "[paperclip-overnight-runner] missing runner script"; exit 1; fi'
user=openclaw
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
priority=15
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME="/data",OPENCLAW_HOME="/data",OPENCLAW_STATE_DIR="/data/.openclaw",OPENCLAW_WORKSPACE_DIR="/data/workspace",XDG_CONFIG_HOME="/data/.config"

[program:webchat-turn-watchdog]
command=/bin/bash -lc 'if [ -x /data/workspace/scripts/webchat-turn-watchdog.sh ]; then exec /data/workspace/scripts/webchat-turn-watchdog.sh; else echo "[webchat-turn-watchdog] missing watchdog script"; exit 1; fi'
user=openclaw
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
priority=16
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME="/data",OPENCLAW_HOME="/data",OPENCLAW_STATE_DIR="/data/.openclaw",OPENCLAW_WORKSPACE_DIR="/data/workspace",XDG_CONFIG_HOME="/data/.config"

[program:tailscaled]
command=/bin/bash -lc 'mkdir -p "${TAILSCALE_STATE_DIR:-/data/.tailscale}" && exec /usr/sbin/tailscaled --state=${TAILSCALE_STATE_DIR:-/data/.tailscale}/tailscaled.state --socket=/tmp/tailscaled.sock --tun=userspace-networking'
autostart=true
autorestart=true
priority=10
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:tailscale-up]
command=/bin/bash -lc '/usr/local/bin/tailscale-up.sh'
autostart=true
# tailscale-up is a one-shot config step; it should exit 0 and NOT be restarted.
autorestart=false
exitcodes=0
startsecs=0
startretries=0
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:syncthing]
command=/usr/local/bin/syncthing-start-guard.sh
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
